# mbsync/mbsync.timer
[Unit]
Description=Run mbsync to refresh mails every minute

[Timer]
OnStartupSec=15s
OnUnitActiveSec=2m

[Install]
WantedBy=timers.target
